# dashboard/com.prism.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prism.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/foxy/.claude/skills/prism/dashboard/server.cjs</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PRISM_DIR</key>
        <string>/tmp/.prism</string>
        <key>PORT</key>
        <string>3333</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/foxy/.prism/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/foxy/.prism/logs/dashboard.log</string>
</dict>
</plist>
